# docker/supervisord.conf
; ============================================
; Supervisor 配置文件
; 管理 A-Plan + EasyTier + TTYD
; ============================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; A-Plan 主服务
[program:a-plan]
command=node src/core/master.js
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/a-plan-error.log
stdout_logfile=/app/logs/a-plan.log
environment=EASYTIER_PATH="/usr/local/bin/easytier-core",TTYD_PATH="/usr/local/bin/ttyd",ET_CONFIG_FILE="/app/configs/easytier.json",NODE_ENV="production"

; 注意：EasyTier 和 TTYD 由 A-Plan 内部管理
; 不需要单独配置 supervisor 条目
; 这样可以通过 Web 界面完全控制
